# deployment/99-ceralive-check-usb-devices.rules
# Elgato/Cam Link (vendor 0fd9) hot(un)plug: poke the CeraUI backend to re-scan
# USB2 capture + audio devices. Target the ceralive.service unit's MAIN pid via
# its cgroup — NOT `pkill -f ceralive`, which substring-matched avahi-daemon's
# process title (`avahi-daemon: registering [ceralive.local]`) and killed mDNS.
# --kill-whom=main preserves the old `pkill -o` single-process intent without
# signalling srtla_send/bcrpt in the same cgroup (SIGUSR2 would terminate them).
ENV{DEVTYPE}=="usb_device", ENV{ID_VENDOR_ID}=="0fd9", RUN+="/usr/bin/systemctl kill --kill-whom=main --signal=SIGUSR2 ceralive.service"
